# passivedns.service
[Unit]
Description=passivedns

[Service]
ExecStart=/usr/bin/passivedns
Restart=on-failure

[Install]
WantedBy=multi-user.target